# clr_debug_daemon.socket.in
[Unit]
Description=Clear Linux OS debuginfo daemon

[Socket]
ListenStream=@SOCKET_PATH@
SocketMode=0600

[Install]
WantedBy=sockets.target
