# clr_debug_daemon.service
[Unit]
Description=Clear Linux debuginfo daemon

[Service]
Type=simple
ExecStart=/usr/bin/clr_debug_daemon

[Install]
WantedBy=multi-user.target
