# clr_debug_daemon.socket
[Unit]
Description=Clear Linux OS debuginfo daemon

[Socket]
ListenStream=/run/clr-debug-info
SocketMode=0600

[Install]
WantedBy=sockets.target
